# notifier/noma_initscript.suse.in
#!/bin/sh
# init file for noma
# v2 - Rune Darrud (theflyingcorpse@gmail.com - 2011.05.26
# v1 - Bernd Erk (berk@netways.de) - 19/12/2007
#
# description: the NETWAYS NOtification MAnager
#
### BEGIN INIT INFO
# Provides:          noma
# Required-Start:    nagios
# Required-Stop:     nagios
# Default-Start:     3 5
# Default-Stop:      0 1 2 6
# Short-Description: the NETWAYS NOtification MAnager
# Description:       the NETWAYS NOtification MAnager
### END INIT INFO
#


# Check for missing binaries (stale symlinks should not happen)
# Note: Special treatment of stop for LSB conformance
NOMA_BIN=@NOTIFIERDIR@/noma_daemon.pl
test -x $NOMA_BIN || { echo "$NOMA_BIN not installed"; 
	if [ "$1" = "stop" ]; then exit 0;
	else exit 5; fi; }

# Check for existence of needed config file and read it
# Read config	

ME="NoMa"
RUNUSER="@RUNAS@"
RUNGROUP="@RUNASGRP@"

. /etc/rc.status

# Reset status of this service
rc_reset


case "$1" in
    start)
	echo -n "Starting $ME "
	/sbin/startproc -u $RUNUSER $NOMA_BIN
	rc_status -v
	;;
    stop)
	echo -n "Shutting down $ME "
	/sbin/killproc $NOMA_BIN
	rc_status -v
	;;
    restart)
	$0 stop
	$0 start
	rc_status
	;;
    status)
	echo -n "Checking for service $ME "
	/sbin/checkproc $NOMA_BIN
	rc_status -v
	;;
    *)
	echo "Usage: $0 {start|stop|status|restart}"
	exit 1
	;;
esac
rc_exit
